# etc/local.d/disable_ipv6.start
#!/bin/bash
# Disable IPv6 net bindings
# Regretfully some components have a dependency on IPv6 and we cannot
# remove it completely, but we still want to get rid of the auto-assigned
# IPv6 addresses causing mayhem
# Copyright (c) 2016 gordonb3 <gordon@bosvangennip.nl>
# License: GPL 3.0+
# NO WARRANTY

eval $(ip -6 -o addr | awk '{print "ip -6 addr del "$4" dev "$2";"}')

